# iso/airootfs/etc/systemd/system/sddm.service.d/10-wait-for-bite.conf
# BITE-OS — SDDM must NOT start until the live user `bite` exists, or
# autologin (User=bite) silently fails and you land on the login screen.
# Drop-in goes on sddm.service (the real unit) so systemd always applies it.
[Unit]
Wants=bite-os-live-setup.service
After=bite-os-live-setup.service
